# data/crontab
*/5 * * * * python /var/www/api.gravitydevelopment.net/cnu/cleandb.py
*/5 * * * * python /var/www/api.gravitydevelopment.net/cnu/checkalerts.py
0,30 * * * * python /var/www/api.gravitydevelopment.net/cnu/checkupdates.py
0 4 * * * python /var/www/api.gravitydevelopment.net/cnu/updatemenu.py